# .config/systemd/user/arch-update.timer.d/override.conf
[Unit]
Description=Run arch-update auto check - 6h

[Timer]
OnUnitActiveSec=
OnUnitActiveSec=6h
