# etc/profile.d/editor.sh
#!/bin/bash

export EDITOR=helix
export SUDO_EDITOR=helix
