# etc/profile.d/fcitx5.sh
#!/bin/bash

export QT_IM_MODULES="wayland;fcitx;ibus"
export XMODIFIERS="@im=fcitx"
