# etc/profile.d/native-file-dialog.sh
#!/bin/bash

export GDK_DEBUG=portals
export GTK_USE_PORTAL=1
