# etc/profile.d/sdl2-wayland.sh
#!/bin/bash

export SDL_VIDEODRIVER="wayland,x11"
