# etc/profile.d/userpath.sh
#!/bin/bash

append_path "$(systemd-path user-binaries)"
export PATH
